# proc/netdpr.service
[Unit]
Description=netdpr uwsgi
After=network.target

[Service]
User=
Group=www-data
WorkingDirectory=
ExecStart=uwsgi --ini uwsgi.ini
